# docker/supervisord.conf
; imgbot — supervisord config that runs both the dashboard and the WhatsApp bot.
;
; Logs from each program are streamed straight to the container's stdout
; (captured by `docker logs`). Each program autorestarts on crash. The bot is
; given a longer startsecs because Chromium takes a beat to boot.

[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/var/run/supervisord.pid

[program:dashboard]
command=python3 -m uvicorn imgbot.api.app:app --host 0.0.0.0 --port 8000
directory=/app
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:bot]
command=node /app/whatsapp-bot/bot.js
directory=/app/whatsapp-bot
autostart=true
autorestart=true
; Give Chromium room to boot before supervisord considers a restart loop.
startsecs=15
stopwaitsecs=15
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
